# systemd/.config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Weekly dotfiles backup timer

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
